# trumpmeltdown.service
[Unit]
Description=TrumpMeltdown

[Service]
ExecStart=/usr/local/bin/trumpmeltdown
