# trumpmeltdown.service
[Unit]
Description=TrumpMeltdown

[Service]
ExecStart=/usr/local/bin/trumpmeltdown
EnvironmentFile=/home/ajn0592/servicevariables
